# backend/docker/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=nobody:nogroup

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
command=php-fpm -F
user=root
autostart=true
autorestart=true
priority=5
startretries=3
startsecs=1
stdout_logfile=/var/log/supervisor/php-fpm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

[program:nginx]
command=nginx -g 'daemon off; error_log /dev/stderr warn;'
user=root
autostart=true
autorestart=true
priority=10
startretries=3
startsecs=1
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

[program:app-init]
command=/bin/bash -c "wait-for-db $DB_HOST $DB_PORT -- php artisan migrate --force && php artisan optimize:clear && php artisan optimize"
user=www-data
directory=/var/www/html
autostart=true
autorestart=false
priority=15
startretries=3
startsecs=1
stdout_logfile=/var/log/supervisor/app-init.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/app-init-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

[eventlistener:processes]
command=sh -c 'echo READY && read line && kill -SIGQUIT $PPID'
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL

